# hyrise_rcm.conf
description "Start hyrise dashboard"
start on runlevel [2345]
stop on shutdown

chdir /home/vagrant/hyrise_rcm
exec /usr/bin/python  app.py > server.log
